# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=apache2ctl -D FOREGROUND

[program:websocket]
command=php /var/www/server.php
stdout_logfile=/var/www/websocket.log
